# systemd/dcss-n8n-caddy-drift.service
[Unit]
Description=DCSS n8n Caddy config drift check

[Service]
Type=oneshot
WorkingDirectory=/opt/dcss-n8n
ExecStart=/opt/dcss-n8n/scripts/check-caddy-drift.sh
