# systemd/dcss-n8n-backup.service
[Unit]
Description=DCSS n8n stack backup and verification
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/dcss-n8n
ExecStart=/opt/dcss-n8n/Docker/backup-stack.sh
ExecStartPost=/opt/dcss-n8n/scripts/verify-backups.sh --latest
ExecStartPost=/opt/dcss-n8n/scripts/backup-retention-report.sh
